# cloudflare-proxy/cloudflared/com.beaconlm.tunnel.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.beaconlab.tunnel</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/REPLACE_ME/.local/bin/cloudflared</string>
    <string>tunnel</string>
    <string>--config</string>
    <string>/Users/REPLACE_ME/path/to/cloudflare-proxy/cloudflared/config.yml</string>
    <string>run</string>
  </array>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/cloudflared-beaconlab.out.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/cloudflared-beaconlab.err.log</string>
</dict>
</plist>
